# backup.service
[Unit]
Description=Backup

[Service]
ExecStart=/home/bread/bin/backup.sh
